# etc/udev/rules.d/99-sdcard.rules
ACTION=="add", KERNEL=="sdf*", RUN+="/usr/local/bin/mv_pics.sh"
